# task8/hello.service.j2
[Unit]
Description=Hello Go application

[Service]
ExecStart=/home/jmaradiaga/DevOps/DevOpsGYMBatch3/hello

[Install]
WantedBy=multi-user.target

